# supervisord.conf
[supervisord]
nodaemon=true
user=user
logfile=/var/log/supervisor/supervisord.log
pidfile=/tmp/supervisord.pid

[program:web]
command=python -m uvicorn app.main:app --host 0.0.0.0 --port 7860
directory=/home/user/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/web.log
stderr_logfile=/var/log/supervisor/web.log
environment=HOME="/home/user"

[program:hermes-gateway]
command=bash /home/user/app/scripts/hermes_gateway.sh
directory=/home/user/app
autostart=true
autorestart=true
stdout_logfile=/var/log/supervisor/gateway.log
stderr_logfile=/var/log/supervisor/gateway.log
environment=HOME="/home/user"
